# etc/profile.d/custom.sh
export PS1=$'\\[\E]633;A\a\\][\\[\\033[0;93m\\]\\u@\\h:\\[\\033[0;91m\\]\\W\\[\\033[00m\\]]\\$ \\[\E]633;B\a\\]'
